# debian/prerm

if [ -x "/etc/init.d/greenback" ]; then
    invoke-rc.d greenback stop || exit $?
fi

