# debian/prerm
#!/bin/sh

if [ -x "/etc/init.d/flashback" ]; then
    invoke-rc.d flashback stop || exit $?
fi

